# modules/oci_k3s_cluster/user_data/additional_vnics/vnics.timer
[Unit]
Description=Vnics timer

[Timer]
OnCalendar=*:0/15
Unit=vnics.service

[Install]
WantedBy=multi-user.target